# deploy/profx/blunderbus-telegram.service
[Unit]
Description=BlunderBus Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=blunderbus
WorkingDirectory=/opt/blunderbus/repo
EnvironmentFile=/opt/blunderbus/.env
ExecStart=/opt/blunderbus/venv/bin/python /opt/blunderbus/repo/scripts/telegram_bot.py
Restart=always
RestartSec=5
StandardOutput=append:/opt/blunderbus/logs/telegram-bot.log
StandardError=append:/opt/blunderbus/logs/telegram-bot.log

[Install]
WantedBy=multi-user.target
