# deploy/ai-workstation/blunderbus-daily-brief.timer
[Unit]
Description=Run BlunderBus daily brief live delivery daily on AI-Workstation

[Timer]
OnCalendar=*-*-* 06:00:00 America/Chicago
Persistent=true
RandomizedDelaySec=5m
Unit=blunderbus-daily-brief.service

[Install]
WantedBy=timers.target
