# deploy/ai-workstation/jarvis-streamdeck.service
[Unit]
Description=Jarvis Stream Deck button daemon
OnFailure=blunderbus-failure-agent@%n.service
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/jarvis-streamdeck
Restart=on-failure
RestartSec=30
Environment=PYTHONUNBUFFERED=1
Environment=STT_MODEL=large-v3-turbo
Environment=STT_DEVICE=cuda
Environment=STT_COMPUTE_TYPE=float16
Environment=STT_LANGUAGE=en
Environment=STT_ENGINE=canary
Environment=CANARY_STT_URL=http://127.0.0.1:8765/transcribe
Environment=CANARY_MAX_NEW_TOKENS=128
Environment=STT_CANARY_FALLBACK=1
Environment=HERMES_BIN=/home/brian/.local/bin/hermes
Environment=JARVIS_RESPONSE_TARGET=discord:general
Environment=JARVIS_DISCORD_SESSION_KEY=agent:main:discord:group:1477789090240921872:1477766287689122015
Environment=JARVIS_DISCORD_CHAT_ID=1477789090240921872
Environment="JARVIS_DISCORD_CHAT_NAME=Blunderbus / #general"
Environment=JARVIS_DISCORD_USER_ID=1477766287689122015
Environment=JARVIS_DISCORD_USER_NAME=Brian
Environment=JARVIS_DISCORD_GUILD_ID=1477768383146491905
Environment=JARVIS_CHAT_TIMEOUT=300
Environment=JARVIS_SEND_TIMEOUT=60
Environment=JARVIS_VOICE_CHAT_TIMEOUT=120
Environment=JARVIS_VOICE_MAX_CHARS=700
Environment=JARVIS_VOICE_SESSION_FILE=/home/brian/voice-notes/ask-voice-session-id.txt
Environment=JARVIS_TTS_ENGINE=edge
Environment=JARVIS_EDGE_TTS_BIN=/home/brian/ai-tools/streamdeck/.venv/bin/edge-tts
Environment=JARVIS_EDGE_VOICE=en-US-GuyNeural
Environment=JARVIS_EDGE_RATE=+5%
Environment=JARVIS_EDGE_PITCH=+0Hz
Environment=JARVIS_ESPEAK_VOICE=en-us
Environment=JARVIS_ESPEAK_SPEED=170

[Install]
WantedBy=default.target
